# build/build-files/DEBIAN/postinst
#!/bin/bash
set -e
systemctl daemon-reload
systemctl enable armController.service
systemctl restart armController.service || true