# build/build-files/DEBIAN/postrm
#!/bin/bash
set -e

case "$1" in
  purge)
    systemctl stop armController.service || true
    systemctl disable armController.service || true
    rm -rf /opt/armController
    ;;
esac